# AX3600/files/etc/uci-defaults/999-QOL_config
uci -q batch << EOI
set pbuf.opt.memory_profile='1gb'
EOI
